# weather_logger.timer
[Unit]
Description=Run weather logger periodically

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
